# razer-config
#!/sbin/openrc-run

depend() {
	need razerd
}

start() {
	if [ -z "${frequency}" ] || [ -z "${resolution}" ]; then
		eerror "You need to set \"frequency\" and \"resolution\" variables"
		return 1
	fi

	ebegin "Adjusting Razer mouse settings"
	/usr/local/bin/razer.py \
		--set-frequency ${frequency} \
		--set-resolution ${resolution}
	eend $?
}
